# package/agent/deb/postinst
#!/bin/sh
set -e

systemctl daemon-reload

case "$1" in
    install)
	systemctl enable clearview-agent
	systemctl start clearview-agent
	;;
esac

