# package/server/deb/postinst
#!/bin/sh
set -e

systemctl daemon-reload

case "$1" in
    install)
	systemctl enable clearview-server
	systemctl start clearview-server
	;;
esac

